# .config/autostart/fix-keyboard.desktop
[Desktop Entry]
Type=Application
Name=fix-keyboard
Exec=xmodmap /home/lexruee/.Xmodmap 

